# filebrowser/roles/filebrowser/templates/filebrowser.service.j2
[Unit]
Description=Filebrowser

[Service]
ExecStart=filebrowser -r /mnt/share
KillMode=process
Restart=always
Type=simple

[Install]
WantedBy=multi-user.target
